# contrib/launchd/com.abevz.af-coordinatord.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.abevz.af-coordinatord</string>

  <key>ProgramArguments</key>
  <array>
    <string>@HOME@/.local/bin/af-coordinatord</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>@HOME@/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>@HOME@</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>@HOME@/Library/Logs/af-coordinatord.log</string>

  <key>StandardErrorPath</key>
  <string>@HOME@/Library/Logs/af-coordinatord.err.log</string>
</dict>
</plist>
